# assets/systemd/transmission.service.d/custom.conf
[Service]
User=root
Environment="EMAIL=$EMAIL"
ExecStart=
ExecStart=/usr/bin/transmission-daemon -f -g /data/transmission
