# assets/systemd/transmission.service.d/custom.conf
[Service]
User=root
ExecStart=
ExecStart=/usr/bin/transmission-daemon -f -g /data/transmission
